# configs/systemd/cmux-vnc.service
[Unit]
Description=Cmux VNC and websockify service
After=network.target
Requires=network.target

[Service]
Type=simple
EnvironmentFile=-/opt/app/app.env
ExecStart=/usr/local/lib/cmux/cmux-rootfs-exec /usr/local/lib/cmux/cmux-start-vnc
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=cmux.target
