# configs/systemd/cmux-vnc.service
[Unit]
Description=Cmux VNC and websockify service
After=network.target
Requires=network.target

[Service]
Type=simple
Environment=DISPLAY=:1
Environment=VNC_PORT=5901
Environment=WEBSOCKIFY_PORT=39380
ExecStart=/usr/local/lib/cmux/cmux-start-vnc
Restart=always
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=cmux.target
