# configs/systemd/cmux-websockify.service
[Unit]
Description=Cmux VNC WebSocket proxy (Rust)
After=cmux-x11vnc.service
Requires=cmux-x11vnc.service

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/log/cmux
ExecStart=/usr/local/bin/vnc-websocket-proxy --listen 0.0.0.0:39380 --target 127.0.0.1:5901 --web /usr/share/novnc
Restart=always
RestartSec=3
StandardOutput=append:/var/log/cmux/vnc-websocket-proxy.log
StandardError=append:/var/log/cmux/vnc-websocket-proxy.log

[Install]
WantedBy=cmux.target
